# debian/systemd/genesis-agent-secrets.service
[Unit]
Description=Prepare NATS runtime secrets and config
Before=genesis-agent-nats.service
After=network.target

[Service]
Type=oneshot
User=genagent
Group=genagent

RuntimeDirectory=genesis_agent
RuntimeDirectoryMode=0700

EnvironmentFile=/etc/genesis_agent/config.env
ExecStart=/opt/genesis_agent/bin/runtime_init
RemainAfterExit=yes

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=genesis-agent-nats.service
